# vieterp/worker/cptuanhuy/supervisord.conf
[supervisord]
nodaemon=true

[program:odoo]
user = odoo
environment = USER="odoo", LOGNAME="odoo", HOME="/var/lib/odoo"
command=/opt/odoo/odoo/odoo-bin -c /etc/odoo/odoo.conf --proxy-mode